# examples/direktor.conf
; Direktor - Example Configuration File
; This format is compatible with Supervisor's INI-style configuration.
;
; Place this file at /etc/direktor/direktor.conf or specify
; an alternative path with: direktord -c /path/to/file.conf

[direktord]
; Path to the supervisor's own log file
logfile = /var/log/direktor/direktord.log
; Log verbosity: debug, info, warn, error
loglevel = info
; PID file written when the daemon starts
pidfile = /var/run/direktor.pid
; Set to true to run in the foreground (useful in containers)
nodaemon = false
; Minimum file descriptors the OS should have available
minfds = 1024
; Minimum OS processes that should be available
minprocs = 200
; Identifier for this supervisor instance
identifier = direktor

; Unix socket for direktorctl communication
socket_path = /var/run/direktor.sock
; File permissions for the control socket (octal)
socket_mode = 0770

; Built-in HTTP server for the web UI and REST API
http_host = 127.0.0.1
http_port = 9876
; Authentication mode: "token" or "basic" (leave empty to disable)
http_auth = token
http_auth_token = changeme

; =========================================================
; Email Notifications
; =========================================================
; Direktor can send email alerts when processes change state.
; Events are batched in a 5-second window to avoid flooding.

[email]
enabled = true
smtp_host = smtp.example.com
smtp_port = 587
username = direktor@example.com
password = app-password-here
from = direktor@example.com
; Comma-separated list of recipient addresses
recipients = ops-team@example.com, oncall@example.com
; Use TLS for the SMTP connection
use_tls = true
; Which state transitions trigger an email (comma-separated)
; Available: RUNNING, STOPPED, EXITED, FATAL, BACKOFF, STARTING
; If omitted, defaults to: RUNNING, STOPPED, EXITED, FATAL
notify_on = FATAL, STOPPED, RUNNING

; =========================================================
; Programme Definitions
; =========================================================
; Each [program:name] section defines a managed process.
; The name after the colon is used as the programme identifier.

[program:webapp]
command = /usr/bin/gunicorn app:application -b 0.0.0.0:8000
directory = /opt/webapp
user = www-data
autostart = true
autorestart = always
startsecs = 5
startretries = 3
stopsignal = TERM
stopwaitsecs = 10
exitcodes = 0,2
priority = 100
stdout_logfile = /var/log/direktor/webapp-stdout.log
stderr_logfile = /var/log/direktor/webapp-stderr.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups = 10
redirect_stderr = false
environment = DJANGO_SETTINGS_MODULE="myapp.settings",SECRET_KEY="supersecret"

[program:worker]
command = /usr/bin/celery -A tasks worker --loglevel=info
directory = /opt/webapp
user = www-data
autostart = true
; Only restart if the exit code is unexpected
autorestart = on-failure
startsecs = 10
startretries = 5
stopsignal = TERM
stopwaitsecs = 30
exitcodes = 0
priority = 200
stdout_logfile = /var/log/direktor/worker-stdout.log
stderr_logfile = /var/log/direktor/worker-stderr.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 5
environment = C_FORCE_ROOT="1"

[program:scheduler]
command = /usr/bin/celery -A tasks beat --loglevel=info
directory = /opt/webapp
user = www-data
autostart = true
autorestart = always
startsecs = 5
startretries = 3
stopsignal = TERM
stopwaitsecs = 10
priority = 300
stdout_logfile = /var/log/direktor/scheduler-stdout.log
stderr_logfile = /var/log/direktor/scheduler-stderr.log
; Merge stderr into stdout log
redirect_stderr = true

[program:nginx-proxy]
command = /usr/sbin/nginx -g "daemon off;"
autostart = true
autorestart = always
startsecs = 2
stopsignal = QUIT
stopwaitsecs = 5
priority = 50
stdout_logfile = /var/log/direktor/nginx-stdout.log
stderr_logfile = /var/log/direktor/nginx-stderr.log

; =========================================================
; Group Definitions
; =========================================================
; Groups let you manage multiple programmes as a single unit.

[group:webstack]
programs = webapp,worker,scheduler
priority = 100

; =========================================================
; Include Additional Configuration Files
; =========================================================
; Use glob patterns to split config across multiple files.
; Relative paths are resolved from this file's directory.

[include]
files = /etc/direktor/conf.d/*.conf
