# deploy/etc/direktor/conf.d/example-app.conf
; Example managed service definition. Copy and edit as needed.

[program:example-app]
command = /usr/bin/python3 -m http.server 8080
directory = /srv/example-app
autostart = true
autorestart = always
startsecs = 2
startretries = 3
stopsignal = TERM
stopwaitsecs = 10
priority = 100
stdout_logfile = /var/log/direktor/example-app-stdout.log
stderr_logfile = /var/log/direktor/example-app-stderr.log
